# supervisor.d/move.ini
[program:move]
command=/bin/move.sh
stdout_logfile=/var/log/move.log
